# src/supervisord.conf
[supervisord]
nodaemon=true

[program:telegram-sonarr-bot]
autostart=true
autorestart=true
user=nobody
command=/usr/bin/node /app/telegram-sonarr-bot-master/sonarr.js